# .provision/udev/50-vagrant-mount.rules
SUBSYSTEM=="bdi", ACTION=="add", RUN+="/bin/bash -c 'sleep 5; /bin/bash /vagrant/.provision/autorun.sh'"
