# .provision/udev/50-vagrant-mount.rules
SUBSYSTEM=="bdi", ACTION=="add", RUN+="/usr/bin/screen -d -m bash -c 'sleep 5; /bin/bash /vagrant/.provision/autorun.sh'"
