# deploy/udev/50-vagrant-mount.rules
SUBSYSTEM=="bdi", ACTION=="add", RUN+="/opt/deploy/udev.sh"
SUBSYSTEM=="bdi", ACTION=="remove", RUN+="/opt/deploy/udev.sh"
